# sys-firmware/pinebookpro-misc/files/openrc-sdio-hciattach
#!/sbin/openrc-run

description="Start bluetooth on ttyS0"

depend()
{
	provide <<sdio-hciattach>>
	need modules devfs
}

start()
{
	ebegin "Start bluetooth on ttyS0"
	/bin/bash -c 'echo 0 > /sys/class/rfkill/rfkill0/state; sleep 1; echo 1 > /sys/class/rfkill/rfkill0/state; sleep 1; exec /usr/bin/hciattach /dev/ttyS0 bcm43xx 1500000 flow nosleep'
	eend $?
}
